# ansible/roles/server-users/templates/reboot-on-pfd.service.j2
[Unit]
Description=Reboot after power failure

[Service]
Type=oneshot

# reboot register for Mac Mini with Intel ISA bridge
ExecStart=setpci -s 0:1f.0 0xa4.b=0

[Install]
WantedBy=sysinit.target
